# lvd/debian/postrm
#!/bin/sh

start(){
    # systemctl disable lvd-detect.service
    systemctl daemon-reload
}

start

exit 0