# lvd/debian/postrm
#!/bin/sh

start(){
    #systemctl disable lvd-detect.service
    systemctl daemon-reload
    rm -rf /opt/ed-lvd/
}

start

exit 0
